# templates/suricata.service.epp
[Unit]
Description=Suricata Intrusion Detection Service
After=syslog.target network-online.target

[Service]
ExecStart=/sbin/suricata -c <%= "${::suricata::config_dir}/suricata.yaml" %> -i eth0 --user <%= $::suricata::user %>

[Install]
WantedBy=multi-user.target

